# tascore/conf/xdg/autostart/qttasserver.desktop
[Desktop Entry]
Exec=/usr/bin/qttasserver
X-Moblin-Priority=High
OnlyShowIn=X-MEEGO-HS;
