# qttasserver/contrib/qttas-server.service
[Unit]
Description=QTTAS Service

[Service]
ExecStart=/usr/bin/qttasserver

[Install]
WantedBy=multi-user.target
