# xsession/qttas-testing.conf
start on BME_OK
console none
script
  /usr/bin/aegis-test-session --set QT_LOAD_TESTABILITY=1
end script